# deploy/codesync-web.service
[Unit]
Description=CodeSync Web UI
After=network.target codesync-server.service
Wants=codesync-server.service

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=-__PROJECT_DIR__/.env
ExecStart=__PROJECT_DIR__/scripts/systemd-web-start.sh
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
