# deploy/codesync-git-sync.timer
[Unit]
Description=Run CodeSync Git Sync Every 5 Minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=codesync-git-sync.service

[Install]
WantedBy=timers.target
