# deploy/codesync-worker.service
[Unit]
Description=CodeSync Worker (job executor)
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=-__PROJECT_DIR__/.env
ExecStart=__PROJECT_DIR__/scripts/systemd-worker-start.sh
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
