# deploy/codesync-git-sync.service
[Unit]
Description=CodeSync Git Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/sixbox/Dev/codesync
ExecStart=/home/sixbox/Dev/codesync/scripts/systemd-git-sync.sh
StandardOutput=journal
StandardError=journal
